# app/configs/pietunes.conf
description "uWSGI server instance configured to serve pitunes"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

chdir /pietunes
exec uwsgi --ini pietunes.ini
